# rootfs/etc/services.d/openssh-server/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid git /usr/sbin/sshd -D -e -p 2222